# bus/dbus.socket.in
[Unit]
Description=D-Bus System Message Bus Socket
Before=sockets.target
ConditionPathExists=!/dev/kdbus/0-system/bus

[Socket]
ListenStream=@DBUS_SYSTEM_SOCKET@
SmackLabelIPIn=*
SmackLabelIPOut=@
