# packaging/dbus.manifest
<manifest>
  <define>
    <domain name="dbus"/>
  </define>
  <request>
    <domain name="dbus"/>
  </request>
  <assign>
    <filesystem path="/usr/bin/dbus-launch" label="_" exec_label="none" />
    <filesystem path="/etc/rc.d/init.d/*" label="_" exec_label="none" />
    <filesystem path="/etc/rc.d/rc3.d/*" label="_" exec_label="none" />
    <filesystem path="/etc/rc.d/rc4.d/*" label="_" exec_label="none" />
  </assign>
</manifest>
